# linux/system/etc/systemd/system/xrunit@.service
[Unit]
Description=Execute the application in X

[Service]
Type=simple
ExecStart=/home/shiro/bin/runinx %I

[Install]
WantedBy=multi-user.target
